# nginx/init.d-gentoo-nginx
#!/sbin/openrc-run

depend() {
    need net
    use dns netmount
    after sshd
}

start() {
    ebegin "Starting nginx"
    start-stop-daemon --start --exec /usr/local/nginx/sbin/nginx
    eend $?
}

stop() {
    ebegin "Stopping nginx"
    kill `cat /usr/local/nginx/logs/nginx.pid`
    eend $?
}

restart() {
    ebegin "Restarting nginx"
    kill -HUP `cat /usr/local/nginx/logs/nginx.pid`
    eend $?
}
